# .config/systemd/user/arxiv-downloader.service
[Unit]
Description=arXiv downloader daemon service

[Service]
Type=simple
Restart=on-failure
RestartSec=30
ExecStart=/usr/bin/env arxiv-downloader
