# .config/systemd/user/autorotate.service
[Unit]
Description=Monitor autorotation service

[Service]
ExecStart=/usr/bin/env autorotate
Restart=always

[Install]
WantedBy=ready.target
