# .config/systemd/user/bcm4331-unhang.service
[Unit]
Description=Service to unhang the BCM4331 chip that has a bad linux driver

[Service]
Type=oneshot
ExecStart=/usr/bin/env wifi.unhang.BCM4331
