# .config/systemd/user/calendar.service
[Unit]
Description=Calendar updates service

[Service]
Type=simple
Restart=on-failure
RestartSec=100
ExecStart=/usr/bin/env calendar.poll
