# .config/systemd/user/cineonu.timer
[Unit]
Description=Fetch CINEONU notifications every minute

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target
