# .config/systemd/user/cineonu.service
[Unit]
Description=CINEONU notifications

[Service]
Type=oneshot
ExecStart=/usr/bin/sh cineonu-notifications
