# .config/systemd/user/clipboard.service
[Unit]
Description=Clipboard recording daemon

[Service]
ExecStart=/bin/bash clipboard.record
Restart=always

[Install]
WantedBy=ready.target
