# .config/systemd/user/dock.service
[Unit]
Description=Trigger action when docked
After=ready.target

[Service]
Type=oneshot
ExecStart=/usr/bin/sh on-dock
