# .config/systemd/user/environment.target
[Unit]
Description=Environment setup
Requires=default.target
After=default.target
