# .config/systemd/user/environment.target
[Unit]
Description=environment is setup
Requires=default.target
After=default.target
