# .config/systemd/user/geolocation.service
[Unit]
Description=Geolocation updates daemon

[Service]
ExecStart=/usr/bin/env sh geolocation.poll
Restart=always

[Install]
WantedBy=ready.target
