# .config/systemd/user/geolocation.timer
[Unit]
Description=Geolocation updates timer

[Timer]
OnActiveSec=45
OnUnitActiveSec=300

[Install]
WantedBy=ready.target
