# .config/systemd/user/github.mark_as_read.timer
[Unit]
Description=Github mark notifications read timer

[Timer]
OnActiveSec=5400
OnUnitActiveSec=5400

[Install]
WantedBy=ready.target
