# .config/systemd/user/github.notifications.service
[Unit]
Description=Memoizes github notifications

[Service]
Type=oneshot
ExecStart=/usr/bin/sh memoize 0 github.notifications
