# .config/systemd/user/github.mark_as_read.timer
[Unit]
Description=Github mark notifications read timer

[Timer]
OnActiveSec=600
OnUnitActiveSec=600

[Install]
WantedBy=ready.target
