# .config/systemd/user/github.notifications.timer
[Unit]
Description=Fetch github notifications hourly

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target
