# .config/systemd/user/ip.link.stat.cache.service
[Unit]
Description=Cache link stats

[Service]
Type=oneshot
ExecStart=/usr/bin/sh ip.link.stat.cache
