# .config/systemd/user/ipgeolocation.timer
[Unit]
Description=IP Geolocation updates timer

[Timer]
OnActiveSec=35
OnUnitActiveSec=600

[Install]
WantedBy=ready.target
