# .config/systemd/user/mbv.timer
[Unit]
Description=Mobile Vikings balance timer

[Timer]
OnActiveSec=90
OnUnitActiveSec=1800

[Install]
WantedBy=ready.target
