# .config/systemd/user/mbv.service
[Unit]
Description=Mobile Vikings daemon

[Service]
ExecStart=/usr/bin/env sh mbv.poll
Restart=always

[Install]
WantedBy=ready.target
