# .config/systemd/user/news.service
[Unit]
Description=News updates daemon

[Service]
ExecStart=/usr/bin/env sh news.poll
Restart=always

[Install]
WantedBy=environment.target
