# .config/systemd/user/papercut.service
[Unit]
Description=Papercut daemon

[Service]
ExecStart=/usr/bin/env sh papercut
Restart=always

[Install]
WantedBy=environment.target
