# .config/systemd/user/ready.timer
[Unit]
Description=Trigger ready target after 10 seconds

[Timer]
Unit=ready.target
OnActiveSec=20
