# .config/systemd/user/ready.target
[Unit]
Description=user instance ready
Requires=environment.target
After=environment.target
