# .config/systemd/user/rtm.service
[Unit]
Description=RTM updates service

[Service]
Type=simple
Restart=on-failure
RestartSec=10
ExecStart=/usr/bin/env memoize-ft 0 rtm.fetch
