# .config/systemd/user/rtm.timer
[Unit]
Description=RTM updates timer

[Timer]
OnActiveSec=60
OnUnitActiveSec=30

[Install]
WantedBy=ready.target
