# .config/systemd/user/slack-notifications@.timer
[Unit]
Description=Slack unread notifications for %I channels timer

[Timer]
OnActiveSec=120
OnUnitActiveSec=120

[Install]
WantedBy=ready.target
