# .config/systemd/user/stib.service
[Unit]
Description=Stib updates daemon

[Service]
ExecStart=/usr/bin/env sh stib.poll
Restart=always

[Install]
WantedBy=ready.target
