# .config/systemd/user/tablet.target
[Unit]
Description=Tablet mode
Requires=ready.target
After=ready.target
