# .config/systemd/user/tv-indexer.service
[Unit]
Description=index tv files in the background

[Service]
ExecStart=/usr/bin/env tv-indexer
Restart=always

[Install]
WantedBy=ready.target
