# .config/systemd/user/ulb.service
[Unit]
Description=ULB service

[Service]
Type=oneshot
ExecStart=/usr/bin/env ulb.poll
