# .config/systemd/user/ulb.service
[Unit]
Description=ULB daemon

[Service]
ExecStart=/usr/bin/env sh ulb.poll
Restart=always

[Install]
WantedBy=environment.target
