# .config/systemd/user/ulb.timer
[Unit]
Description=ULB timer

[Timer]
OnActiveSec=5
OnUnitActiveSec=5

[Install]
WantedBy=ready.target
