# .config/systemd/user/undock.service
[Unit]
Description=Trigger action when undocked
After=ready.target

[Service]
Type=oneshot
ExecStart=/usr/bin/sh on-undock
