# .config/systemd/user/urlab.service
[Unit]
Description=UrLaB updates service

[Service]
Type=simple
Restart=on-failure
RestartSec=100
ExecStart=/usr/bin/sh memoize-ft 0 urlab.fetch
