# .config/systemd/user/urlab.service
[Unit]
Description=UrLaB updates daemon

[Service]
ExecStart=/usr/bin/env sh urlab.poll
Restart=always

[Install]
WantedBy=default.target

