# .config/systemd/user/voo.timer
[Unit]
Description=VOO timer

[Timer]
OnActiveSec=150
OnUnitActiveSec=1800

[Install]
WantedBy=ready.target
