# .config/systemd/user/voo.service
[Unit]
Description=VOO daemon

[Service]
ExecStart=/usr/bin/env sh voo.poll
Restart=always

[Install]
WantedBy=default.target
