# .config/systemd/user/vulnerabilities.service
[Unit]
Description=Memoizes audit's output

[Service]
Type=simple
Restart=on-failure
RestartSec=300
ExecStart=/usr/bin/env memoize-ft 0 audit
