# .config/systemd/user/vulnerabilities.service
[Unit]
Description=Memoizes arch-audit's output

[Service]
Type=oneshot
ExecStart=/usr/bin/sh memoize 0 arch-audit -qu
