# .config/systemd/user/weather.service
[Unit]
Description=Weather updates daemon

[Service]
ExecStart=/usr/bin/env sh weather.poll
Restart=always

[Install]
WantedBy=ready.target

