# .config/systemd/user/wifi-fingerprint.timer
[Unit]
Description=Wifi fingerprinting timer

[Timer]
OnActiveSec=35
OnUnitActiveSec=600

[Install]
WantedBy=ready.target
